# aethermapd.service
[Unit]
Description=Aethermap Input Remapping Daemon
After=network.target
Wants=udev-settle.target

[Service]
Type=simple
ExecStart=/usr/bin/aethermapd
Restart=on-failure
RestartSec=5

# Socket location - systemd creates /run/aethermap/ for us
RuntimeDirectory=aethermap
RuntimeDirectoryMode=0755

# Ensure config directory exists
ExecStartPre=/usr/bin/mkdir -p /etc/aethermap

# Environment
Environment=RUST_LOG=info

# Hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/aethermap /run/aethermap
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
